# scripts/services/nsqd.service
[Unit]
Description=NSQDaemon

Requires=docker.service
Requires=etcd.service

After=docker.service
After=etcd.service

[Service]
TimeoutStartSec=0

EnvironmentFile=/etc/environment

ExecStartPre=-/usr/bin/docker kill nsqd
ExecStartPre=-/usr/bin/docker rm nsqd
ExecStartPre=/usr/bin/docker pull dockerfile/nsq
ExecStartPre=/usr/bin/docker pull kzvezdarov/ircdd-nsq-confd

ExecStart=/usr/bin/docker run \
                          -e HOST_IP=${COREOS_PUBLIC_IPV4} \
                          -p 4150:4150 -p 4151:4151 \
                          --name nsqd \
                          kzvezdarov/ircdd-nsq-confd \
                          confd-watch-nsqd

ExecStop=/usr/bin/docker stop nsqd

[X-Fleet]
Global=true
